# admin/deploy/fbrk-enrich.service
[Unit]
Description=FBRK AI enrichment — process new articles
After=network-online.target fbrk-admin.service
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/fbrk-admin
EnvironmentFile=/opt/fbrk-admin/.env
ExecStart=/opt/fbrk-admin/.venv/bin/python3 /opt/fbrk-admin/enrich.py --sleep 15
StandardOutput=append:/opt/fbrk-admin/enrich-cron.log
StandardError=append:/opt/fbrk-admin/enrich-cron.log
TimeoutStartSec=3600
